# supervisord.conf
[supervisord]
nodaemon=true

[program:service]
command=/app/service/service